# players/vuln_box/systemd/log.service
[Unit]
Description=Log

[Service]
User=root
Type=simple
ExecStart=/usr/bin/python /root/log.py

[Install]
WantedBy=multi-user.target
