# players/vuln_box/systemd/tunn.service
[Unit]
Description=Tunn

[Service]
User=tunnel
Type=simple
ExecStart=/usr/bin/python /home/tunnel/tunn.py

[Install]
WantedBy=multi-user.target
